# init/supervisor.conf
[supervisord]
nodaemon = true

[program:nzbget]
user = nobody
group = users
autorestart = False
command = nzbget -s -c /config/nzbget.conf
environment = TERM=xterm
